# var/service/autofs/run
#!/bin/sh
# $Header:  Exp $

exec automount -f

