# var/service/ntpd/run
#!/bin/sh

# net-misc/openntpd supervise script

sleep 10

rm -f /run/ntpd.sock
exec /usr/sbin/ntpd -n

